# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid
user=root

[unix_http_server]
file=/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

; ─── NGINX Unit (web server) ────────────────────────────────────────────────
[program:unit]
command=unitd --no-daemon --log /dev/stderr
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Laravel Reverb (WebSocket server) ──────────────────────────────────────
[program:reverb]
command=php artisan reverb:start --host=0.0.0.0 --port=8080 --no-interaction
directory=/var/www/html
user=unit
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Laravel Queue Worker ────────────────────────────────────────────────────
[program:queue]
command=php artisan queue:work --sleep=3 --tries=3 --max-time=3600 --no-interaction
directory=/var/www/html
user=unit
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Laravel Scheduler ───────────────────────────────────────────────────────
[program:scheduler]
command=php artisan schedule:work --no-interaction
directory=/var/www/html
user=unit
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
